# tomcat7/centos/src/tomcat.sv.conf
[program:tomcat7]
command=/opt/launch_tomcat.sh
redirect_stderr=true
